# packages/base/busybox/fs-target/etc/inittab
null::sysinit:/etc/rc.d/rcS boot

null::respawn:/sbin/syslogd -n -K -t
::respawn:-/bin/sh

null::restart:/sbin/init
null::ctrlaltdel:/sbin/reboot
null::shutdown:/etc/rc.d/rcS shutdown
